# etc/systemd/system/ofono.service.d/ril-override.conf
[Service]
Environment=OFONO_RIL_DEVICE=ril
Environment=OFONO_RIL_NUM_SIM_SLOTS=1
Environment=OFONO_RIL_RAT_LTE=1
ExecStart=
ExecStart=/usr/sbin/ofonod --nodetach

[Install]
After=lxc@android.service
