# etc/profile.d/webengine-workaround.sh
export QTWEBENGINE_DISABLE_GPU_THREAD=1
